# whisprd.service
[Unit]
Description=Real-time Whisper Dictation System (Whisprd)
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
User=%i
Group=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/%i
WorkingDirectory=/home/%i
ExecStart=/usr/bin/python3 /path/to/whisprd/whisprd_cli.py --daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/%i/.config /home/%i/.cache /tmp

# Audio and input permissions
SupplementaryGroups=audio input

[Install]
WantedBy=default.target 